# cron-job-with-docker/cron-tab-external/crontab
SHELL=/bin/bash
* * * * * root /usr/local/bin/find-me.sh >> /var/log/cron.log 2>&1
